# docs/pdf-parser@.service
[Unit]
Description=PDF Parser API (%i)
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/code/fin_parser
Environment=PYTHONUNBUFFERED=1
Environment=PATH=/home/%i/.local/bin:/usr/local/bin:/usr/bin
ExecStart=/home/%i/.local/bin/uv run pdf-parser serve --host 0.0.0.0 --port 8005
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
